# supervisord.conf
[supervisord]
nodaemon=true

[program:webssh]
command=service ssh start